# etc/xdg/autostart/crispy.desktop
[Desktop Entry]
Name=Crispy
Comment=
TryExec=crispy
Exec=crispy
Type=Application
OnlyShowIn=XFCE;
Terminal=false
Hidden=false
